# package/midi-link/midi-link.service
[Unit]
Description=MIDI Link Application

[Service]
ExecStart=/usr/bin/MidiLink
Restart=always

[Install]
WantedBy=multi-user.target

